# files/inittab
# /etc/inittab
::sysinit:/hab/bin/hab pkg exec core/coreutils echo STARTING SYSTEM
::sysinit:/etc/init.d/startup
:3:initdefault:
:3:respawn:/hab/bin/hab pkg exec core/busybox-static runsvdir /etc/rc.d
tty1::respawn:/hab/bin/hab pkg exec core/busybox-static getty 38400 tty1
tty2::respawn:/hab/bin/hab pkg exec core/busybox-static getty 38400 tty2
tty3::respawn:/hab/bin/hab pkg exec core/busybox-static getty 38400 tty3
::ctrlaltdel:/hab/bin/hab pkg exec core/util-linux umount -a -r
::shutdown:/hab/bin/hab pkg exec core/coreutils echo SHUTTING DOWN
::shutdown:/hab/bin/hab pkg exec core/util-linux umount -a -r
# end of /etc/inittab
